# site/minecraft/templates/minecraft.service.epp
[Unit]
Description=Minecraft Server

Wants=networt.target
After=network.target

[Service]
WorkingDirectory=<%= $install_dir %>
ExecStart=/usr/bin/java -Xmx512m -Xms32m -jar server.jar nogui

[Install]
WantedBy=multi-user.target
